# launchd/com.simp.cowork_bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.simp.cowork_bridge</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3.10</string>
        <string>REPLACE_WITH_REPO_PATH/bin/cowork_bridge.py</string>
        <string>--port</string>
        <string>8767</string>
        <string>--repo-path</string>
        <string>REPLACE_WITH_REPO_PATH</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>CLAUDE_CODE_CMD</key>
        <string>claude</string>
        <key>COWORK_TIMEOUT</key>
        <string>300</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>REPLACE_WITH_HOME/bullbear/logs/cowork_bridge.log</string>
    <key>StandardErrorPath</key>
    <string>REPLACE_WITH_HOME/bullbear/logs/cowork_bridge_err.log</string>
    <key>WorkingDirectory</key>
    <string>REPLACE_WITH_REPO_PATH</string>
</dict>
</plist>
